# example.service
[Unit]
Description=Pirate pal service

[Service]
ExecStart=/opt/pirate-pal/pirate-pal

[Install]
WantedBy=multi-user.target
